# hal-tts.service
[Unit]
Description=HAL TTS Streaming Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=tts
Group=tts
WorkingDirectory=/opt/tts-pipeline
ExecStart=/opt/tts-pipeline/streaming_safe_daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hal-tts

# Performance optimizations for audio processing
LimitNOFILE=65536
LimitNPROC=4096
Nice=-5

# Security settings (relaxed for audio access)
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/opt/tts-pipeline /tmp /var/log /run/user/1000
PrivateTmp=false
PrivateDevices=false
ProtectControlGroups=false
ProtectKernelModules=false
ProtectKernelTunables=false
RestrictRealtime=false
SupplementaryGroups=audio

# Environment variables
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="HOME=/opt/tts-pipeline"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="PULSE_RUNTIME_PATH=/run/user/1000/pulse"
Environment="HAL_DOTMATRIX_ENABLED=1"
Environment="HAL_DOTMATRIX_QUEUE_DIR=/tmp/halo-dotmatrix/queue"
Environment="HAL_DOTMATRIX_WAV_DIR=/tmp/halo-dotmatrix/wav"

[Install]
WantedBy=multi-user.target